# .docker/app/crontabs/magento.cron
* * * * * cd /app && php bin/magento cron:run 2>&1 | tee -a /var/log/supervisor/magento-cron.log
